# install_scripts/templates/upstart/replicated-operator-stop.conf
description "Replicated Operator shutdown script"
author "Replicated.com"
start on stopping replicated-operator
kill timeout 30
script
    exec /usr/bin/docker stop replicated-operator
end script
