# install_scripts/templates/replicated-init-aws-bootstrap.conf
description "Replicated Bootstrap"
author "Replicated.com"

start on filesystem and started docker and runlevel [2345]
stop on runlevel [!2345]

exec /etc/replicated-bootstrap/init-defaults.sh